# superconf/mongo.conf
[program:mongo]
command=/usr/bin/mongod --dbpath /mongodb/data
autorestart=false
