# superconf/server.conf
[program:server]
directory=/root/randomery
command=/bin/bash -c ". %(directory)/venv/bin/activate && uwsgi -H %(directory)/venv --ini uwsgi.ini"
autorestart=false
