# superconf/worker.conf
[program:worker]
directory=/root/randomery
command=/bin/bash -c ". %(directory)/venv/bin/activate && python -c 'from lib.worker import job_loop;job_loop()'"
autorestart=false
